# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=php-fpm -F
priority=10
autorestart=true

[program:nginx]
command=nginx -g "daemon off;"
priority=20
autorestart=true